# supabase/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reexec
systemctl daemon-reload

# kong
systemctl enable kong
systemctl start kong